# deploy/systemd/sniper-publisher.service
[Unit]
Description=SniperBO Official Dashboard Publisher
After=network-online.target sniper-signals-api.service
Wants=network-online.target
Requires=sniper-signals-api.service

[Service]
Type=simple
User=sniper
Group=sniper
WorkingDirectory=/opt/iasniperbo
EnvironmentFile=/opt/iasniperbo/scripts/official_publisher.local.env
ExecStart=/opt/iasniperbo/.venv/bin/python /opt/iasniperbo/scripts/official_dashboard_publisher.py --env-file /opt/iasniperbo/scripts/official_publisher.local.env --local-url http://127.0.0.1:8787/dashboard --interval 0.7 --log-file /opt/iasniperbo/official_dashboard_publisher.log
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
